# ansible/templates/cloudflared.service.j2
# Cloudflared DNS-over-HTTPS Proxy Service
# Generated by Ansible for {{ inventory_hostname }}
# Profile: {{ profile }}

[Unit]
Description=Cloudflared DNS-over-HTTPS Proxy
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-apps/run-tunnel/run-as-service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/cloudflared --config /etc/cloudflared/config.yml
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/cloudflared

[Install]
WantedBy=multi-user.target
